# jsonserver-watch.path
# /etc/systemd/system/jsonserver-watch.path
[Unit]
Description=Watch for CORS server file changes

[Path]
PathChanged=/usr/local/bin/jsonserver.py
Unit=jsonserver-restart.service

[Install]
WantedBy=multi-user.target